# systemd/system/carbon-pool.service
[Unit]
Description=carbon-pool daemon

[Service]
ExecStart=/usr/bin/carbon-pool
User=daemon
Group=daemon

[Install]
WantedBy=multi-user.target
